# mysql/backup.cron
*/1 * * * * echo $(date +%F%n%T) >> /test.txt
